# packaging/debian/circus-server.service
[Unit]
Description=Circus daemon
After=network.target

[Service]
Type=simple
ExecStart=/usr/lib/circus/server.exe
Restart=on-failure
